# api/tflops-api.service
[Unit]
Description=TFLOPS Global Intelligence API
After=network.target postgresql.service

[Service]
Type=simple
User=root
# Ensure this directory exists on the server
WorkingDirectory=/var/www/tflops-api
# Load environment variables from a secure file
EnvironmentFile=/etc/tflops.env
ExecStart=/var/www/tflops-api/tflops-api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
